# roles/amzn-2/templates/etc/systemd/system/partition.mount.j2
[Mount]
What=/dev/disk/by-uuid/{{ blkid }}
Where={{ mount_path }}
Type=ext4
Options={{ mount_opts }}
[Install]
WantedBy=multi-user.target